# build_scripts/supervisor.filebeat.conf
[program:filebeat]
command=/usr/bin/filebeat -c /etc/filebeat/filebeat.yml
autostart=true
autorestart=true
redirect_stderr=true
